# src/crontabs/ec2_crontab
* * * * * echo "Job1 initiated from Jenkins" >> ~/test/out1.txt
